# files_example/etc/uci-defaults/97-disable-services
#!/bin/sh

# disable docker
/etc/init.d/dockerd stop
/etc/init.d/dockerd disable

exit 0
